# services/ml-analytics-service-modular/ml-analytics-monitoring.service
[Unit]
Description=ML Analytics Continuous Monitoring Service
Documentation=man:systemd.service(5)
After=network.target network-online.target ml-analytics-lxc.service
Wants=network-online.target
Requires=ml-analytics-lxc.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=mdoehler
Group=mdoehler
WorkingDirectory=/home/mdoehler/aktienanalyse-ökosystem/services/ml-analytics-service-modular
ExecStart=/usr/bin/python3 continuous_monitoring_service_v1_0_0_20250821.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=30
KillMode=mixed
TimeoutStopSec=20

# Environment
Environment=PYTHONPATH=/home/mdoehler/aktienanalyse-ökosystem/services/ml-analytics-service-modular
Environment=MONITORING_INTERVAL=30
Environment=LXC_CONTAINER_IP=10.1.1.174

# Security & Resource Limits
NoNewPrivileges=true

# Resource Limits
LimitNOFILE=1024
LimitNPROC=512
MemoryAccounting=true
MemoryMax=256M
CPUAccounting=true
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ml-analytics-monitoring

[Install]
WantedBy=multi-user.target