# deployment/systemd/ml-scheduler.service
[Unit]
Description=ML Scheduler Service - Geplantes Model Training
Documentation=file:///home/aktienanalyse/aktienanalyse-ökosystem/README.md
After=network.target postgresql.service redis.service ml-analytics.service
Wants=network.target
Requires=postgresql.service redis.service

[Service]
Type=oneshot
User=aktienanalyse
Group=aktienanalyse
WorkingDirectory=/home/aktienanalyse/aktienanalyse-ökosystem/services/ml-analytics-service-modular

# Scheduler Script ausführen
ExecStart=/home/aktienanalyse/aktienanalyse-ökosystem/deployment/scripts/ml-training-scheduler.sh

# Environment Variables
Environment=PYTHONPATH=/home/aktienanalyse/aktienanalyse-ökosystem/services/ml-analytics-service-modular
Environment=ML_DATABASE_HOST=localhost
Environment=ML_DATABASE_PORT=5432
Environment=ML_DATABASE_NAME=aktienanalyse
Environment=ML_DATABASE_USER=ml_service
Environment=ML_DATABASE_PASSWORD=ml_service_secure_2025
Environment=ML_REDIS_URL=redis://localhost:6379/2
Environment=ML_LOG_LEVEL=INFO

# Keine automatischen Restarts für oneshot
Restart=no

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/aktienanalyse/aktienanalyse-ökosystem/logs /tmp

# Resource Limits
LimitNOFILE=1024
LimitNPROC=512
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ml-scheduler

[Install]
WantedBy=multi-user.target