# systemd/aktienanalyse-event-bus-modular.service
[Unit]
Description=Aktienanalyse Event Bus Service Modular
After=network.target redis-server.service rabbitmq-server.service
Requires=redis-server.service rabbitmq-server.service

[Service]
Type=simple
User=aktienanalyse
Group=aktienanalyse
WorkingDirectory=/opt/aktienanalyse-ökosystem/services/event-bus-service-modular
ExecStart=/opt/aktienanalyse-ökosystem/services/event-bus-service-modular/start_service.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONPATH=/opt/aktienanalyse-ökosystem:/opt/aktienanalyse-ökosystem/shared
Environment=EVENT_BUS_SERVICE_PORT=8014
Environment=REDIS_URL=redis://localhost:6379/0
Environment=RABBITMQ_URL=amqp://guest:guest@localhost:5672/

[Install]
WantedBy=multi-user.target