# services/event-bus-service/event-bus-service.service
[Unit]
Description=Event-Bus Service v7.0.0 - KRITISCHES HERZSTÜCK Event Hub
Documentation=https://github.com/aktienanalyse-ökosystem/services/event-bus-service
After=network.target redis-server.service
Requires=network.target
Wants=redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aktienanalyse-ökosystem/services/event-bus-service
ExecStart=/opt/aktienanalyse-ökosystem/services/event-bus-service/venv/bin/python event_bus_minimal_v7_0_0_20250824.py
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Environment Variables
Environment=PYTHONPATH=/opt/aktienanalyse-ökosystem
Environment=EVENT_BUS_HOST=0.0.0.0
Environment=EVENT_BUS_PORT=8014
Environment=LOG_LEVEL=INFO
Environment=REDIS_HOST=localhost
Environment=REDIS_PORT=6379
Environment=REDIS_EVENT_DB=0
Environment=POSTGRES_HOST=localhost
Environment=POSTGRES_PORT=5432
Environment=POSTGRES_DB=aktienanalyse_ecosystem
Environment=POSTGRES_USER=aktienanalyse_user
Environment=POSTGRES_PASSWORD=secure_password_2024
Environment=EVENT_CHANNEL_PREFIX=events
Environment=EVENT_TTL_SECONDS=86400
Environment=MAX_EVENT_SIZE=1048576

# Security Settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aktienanalyse-ökosystem/services/event-bus-service
ReadWritePaths=/tmp

# Resource Limits
LimitNOFILE=65536
LimitCORE=0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=event-bus-service

[Install]
WantedBy=multi-user.target
Alias=event-bus.service