# services/market-data-service/market-data-service.service
[Unit]
Description=Market Data Service Enhanced v6.0.0 - Yahoo Finance Integration Service
After=network.target redis-server.service postgresql.service event-bus-service.service
Requires=network.target
Wants=redis-server.service postgresql.service event-bus-service.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aktienanalyse-ökosystem/services/market-data-service
ExecStart=/opt/aktienanalyse-ökosystem/services/market-data-service/venv/bin/python market_data_service_v6_0_0_20250824.py
Restart=always
RestartSec=10
StartLimitBurst=3

# Environment Variables
Environment=PYTHONPATH=/opt/aktienanalyse-ökosystem
Environment=SERVICE_HOST=0.0.0.0
Environment=SERVICE_PORT=8002
Environment=LOG_LEVEL=INFO
Environment=REDIS_HOST=localhost
Environment=REDIS_PORT=6379
Environment=POSTGRES_HOST=localhost
Environment=POSTGRES_PORT=5432
Environment=POSTGRES_DB=aktienanalyse_ecosystem
Environment=POSTGRES_USER=aktienanalyse_user
Environment=POSTGRES_PASSWORD=secure_password_2024
Environment=EVENT_BUS_URL=http://localhost:8014

# Resource Limits
LimitNOFILE=65536
LimitCORE=0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=market-data-service

[Install]
WantedBy=multi-user.target
Alias=market-data.service