# configs/systemd/iex-cloud-microcap.service
[Unit]
Description=IEX Cloud Microcap Data Source Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=aktienanalyse
Group=aktienanalyse
WorkingDirectory=/opt/aktienanalyse-ökosystem/services/data-sources
ExecStart=/usr/bin/python3 iex_cloud_microcap_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/opt/aktienanalyse-ökosystem
Environment=IEX_CLOUD_TOKEN=demo

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aktienanalyse-ökosystem/logs

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target