# configs/systemd/marketstack-global-exchanges.service
[Unit]
Description=Marketstack Global Exchanges Service
Documentation=man:marketstack-global-exchanges(8)
After=network.target

[Service]
Type=simple
User=aktienanalyse
Group=aktienanalyse
WorkingDirectory=/opt/aktienanalyse-ökosystem/services/data-sources
ExecStart=/usr/bin/python3 marketstack_global_exchanges_v1.0.0_20250817.py
Restart=always
RestartSec=12

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/aktienanalyse-ökosystem

# Resource limits
MemoryLimit=1G
CPUQuota=70%
TasksMax=80
LimitNOFILE=65536

# Environment
Environment=PYTHONPATH=/opt/aktienanalyse-ökosystem
Environment=MARKETSTACK_API_KEY=demo_key_marketstack_global

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=marketstack-global-exchanges

[Install]
WantedBy=multi-user.target